# numlockOnTty.service
[Unit]
Description=numlock

[Service]
ExecStart=/usr/bin/numlockOnTty
StandardInput=tty
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
